# services/remoteusb-usbipd.service
[Unit]
Description=RemoteUSB USB/IP Daemon
After=network.target

[Service]
Type=forking
ExecStartPre=/sbin/modprobe usbip_core
ExecStartPre=/sbin/modprobe usbip_host
ExecStart=/usr/sbin/usbipd -D
# Bereits angesteckte Geräte nachholen (udev greift nur bei add-Events).
ExecStartPost=/bin/sh -c 'sleep 1 && /usr/local/bin/remoteusb-usb-autobind.sh || true'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
