# services/remoteusb-watchdog.service
[Unit]
Description=RemoteUSB Watchdog (WLAN + WireGuard Überwachung)
After=network.target remoteusb-gpio.service
Wants=remoteusb-gpio.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/remoteusb-watchdog.py
Restart=always
RestartSec=5
EnvironmentFile=/etc/remoteusb/watchdog.env

[Install]
WantedBy=multi-user.target
